# root/etc/cron.d/rsync
# run rsync
0 2 * * 3 root /config/rsync.sh 2>&1 >> /config/logs/log.txt
